# scripts/tansaeng-daemon-report.cron
# tansaeng 데몬 재시작 일일 리포트 — 매일 저녁 9시(21시) 텔레그램 발송
0 21 * * * root /var/www/html/scripts/daemon_restart_report.sh >> /var/log/tansaeng-daemon-report.log 2>&1
